# cron/crontab
*/10 * * * * root *graphAPI/countAllPlaceSocialCount